# contrib/systemd/gh-traffic.timer
[Unit]
Description=Run gh-traffic daily to accumulate GitHub traffic history

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
